# deploy/example-project.conf
description "API uwsgi upstart"

start on runlevel [2345]
stop on runlevel [!2345]

kill timeout 5
respawn

script
exec /home/lontong/www/example-project/deploy/start.sh
end script

